# scrapy/supervisor.gp.conf
[program:gp]
command=sh run.sh
directory=/app
user=root
autorestart=true
startretries=99999
stdout_logfile=/var/log/supervisor/gp.out.log
stderr_logfile=/var/log/supervisor/gp.error.log